# cv_test/garage/motion/cron.hourly/motion
#!/bin/sh

find /var/lib/motion/ -mtime +6 -type f -delete
